# layout/DEBIAN/postinst
#!/bin/sh
PREFIX=""
if [ -d /var/jb ]; then
    PREFIX="/var/jb"
fi

# Ensure the decoy binary is executable
chmod 755 $PREFIX/usr/libexec/corelliumd 2>/dev/null || true
chown root:wheel $PREFIX/usr/libexec/corelliumd 2>/dev/null || true

chown root:wheel $PREFIX/Library/LaunchDaemons/c.eolnmsuk.corelliumdecoy.plist
chmod 644 $PREFIX/Library/LaunchDaemons/c.eolnmsuk.corelliumdecoy.plist
$PREFIX/usr/bin/launchctl unload $PREFIX/Library/LaunchDaemons/c.eolnmsuk.corelliumdecoy.plist 2>/dev/null || true
$PREFIX/usr/bin/launchctl load $PREFIX/Library/LaunchDaemons/c.eolnmsuk.corelliumdecoy.plist 2>/dev/null || true
exit 0
